# overlay/usr/lib/systemd/system/fs-initialize.service
[Unit]
Description=Initialize Filesystems
ConditionPathExists=!/mnt/overlay
Before=sysinit.target
After=local-fs.target
DefaultDependencies=no

[Service]
Type=oneshot
ExecStart=/sbin/mkfs.ext4 -FFq /dev/mmcblk0p3
ExecStartPost=/bin/systemctl start fsw-rw.mount
ExecStartPost=/bin/mkdir /mnt/tmp /mnt/overlay
ExecStartPost=/bin/systemctl start fs-overlay.mount

[Install]
WantedBy=sysinit.target.wants
